# templates/oauth2_proxy@.service.erb
[Unit]
Description=OAuth2 Proxy

[Service]
User=<%= scope['::oauth2_proxy::user'] %>
Group=<%= scope['::oauth2_proxy::group'] %>
ExecStart=<%= scope['::oauth2_proxy::install_root'] %>/bin/oauth2_proxy --config=/etc/oauth2_proxy/%i.conf

[Install]
WantedBy=multi-user.target
