# templates/oauth2_proxy.service.erb
[Unit]
Description=OAuth2 Proxy

[Service]
User=<%= scope['::oauth2_proxy::user'] %>
Group=<%= scope['::oauth2_proxy::group'] %>
ExecStart=/opt/oauth2_proxy/bin/oauth2_proxy --config=/etc/oauth2_proxy/oauth2_proxy.conf

[Install]
WantedBy=multi-user.target
